# deploy/metadata.service
[Unit]
Description=File Metadata Inspector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=infra
Group=infra
WorkingDirectory=/opt/metadata/app
Environment=NODE_ENV=production
EnvironmentFile=/etc/metadata/metadata.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/metadata/app
MemoryMax=320M

[Install]
WantedBy=multi-user.target
